# systemd/nas-monitor.timer
# Managed by NASe — do not edit manually. Re-run apply.sh instead.
[Unit]
Description=NASe drive health monitor timer (daily SMART check)

[Timer]
# Run at 06:00 every day
OnCalendar=*-*-* 06:00:00
# Also run 15 minutes after each boot (catches issues after a restart)
OnBootSec=15min
Persistent=true
Unit=nas-monitor.service

[Install]
WantedBy=timers.target
